# com.eirinior.greek.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eirinior.greek.digest</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Point to whatever python has anthropic + openai installed -->
        <string>/Users/rena/anaconda3/bin/python3</string>
        <string>/Users/rena/greek-weekly-digest/generate_weekly_digest.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>REPLACE_WITH_YOUR_KEY</string>
        <key>OPENAI_API_KEY</key>
        <string>REPLACE_WITH_YOUR_KEY</string>
        <key>GITHUB_TOKEN</key>
        <string>REPLACE_WITH_YOUR_TOKEN</string>
        <key>PATH</key>
        <string>/Users/rena/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Every Monday at 09:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/rena/greek-weekly-digest</string>

    <key>StandardOutPath</key>
    <string>/Users/rena/greek-weekly-digest/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rena/greek-weekly-digest/logs/stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
